# examples/systemd/emqx_exporter.socket
[Unit]
Description=Node Exporter

[Socket]
ListenStream=8085

[Install]
WantedBy=sockets.target
